# scripts/daily-csv-routine/launchd/com.vo360.csv-routine.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vo360.csv-routine</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Users/urielholzman/ghl-sub-account-builder/scripts/daily-csv-routine/run-daily.mjs</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/urielholzman/ghl-sub-account-builder</string>

    <key>StandardOutPath</key>
    <string>/Users/urielholzman/Desktop/CSV Ready Zips/_logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/urielholzman/Desktop/CSV Ready Zips/_logs/launchd-stderr.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
